# systemd/lan-transcriber-api.service
[Unit]
Description=LAN Transcriber API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
ExecStartPre=/usr/bin/python3 -m lan_app.db_init
ExecStart=/usr/bin/python3 -m uvicorn lan_app.api:app --host 0.0.0.0 --port 7860
Restart=on-failure

[Install]
WantedBy=multi-user.target
