# systemd/recording-transcriber.service
[Unit]
Description=DEPRECATED: LAN Recording-Transcriber wrapper (use lan-transcriber-api.service)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/lan-transcriber
Environment=PYTHONPATH=/opt/lan-transcriber
ExecStart=/usr/bin/python3 -m uvicorn lan_app.api:app --host 0.0.0.0 --port 7860
Restart=on-failure

[Install]
WantedBy=multi-user.target
