# systemd/recording-transcriber.service
[Unit]
Description=DEPRECATED: LAN Recording-Transcriber wrapper (use lan-transcriber-api.service)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/%i/lan-transcriber
Environment=PYTHONPATH=/home/%i/lan-transcriber
Environment=LAN_REDIS_URL=redis://127.0.0.1:6379/0
ExecStart=/home/%i/asr_env/bin/python -m uvicorn lan_app.api:app --host 0.0.0.0 --port 7860
Restart=on-failure

[Install]
WantedBy=multi-user.target
